# systemd/rn-remote-relay.service
[Unit]
Description=RNRemote Relay Server
After=network.target postgresql.service

[Service]
Type=simple
User=remotelink
Group=remotelink
WorkingDirectory=/opt/rn-remote
EnvironmentFile=/etc/rn-remote/env
ExecStart=/opt/rn-remote/venv/bin/python relay-server/relay.py --host 0.0.0.0 --port 8765
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
